# frood/root/etc/local.d/hello.start
#!/bin/sh

(
    echo "Hello, world!"
    echo
    echo "Image: $(cat /etc/frood-release)"
    echo
    uname -a
    echo
    free -h
    echo
    zpool status
    echo
    ip addr
) | mail-alert "frood boot"
