# frood/root/etc/profile.d/motd.sh
echo
zpool status
echo
echo "version: $(cat /etc/frood-release)"
echo
